# ansible/roles/bot/templates/bot.service.j2
[Unit]
Description=TripleA Bot %i
Documentation=

[Service]
WorkingDirectory={{ bot_install_folder }}
User=bot
Group=bot
Environment="MAPS_FOLDER={{ bot_maps_folder }}"
ExecStart=java -server \
  -Xmx{{ bot_max_memory }} \
  -Xss1250K \
  -Djava.awt.headless=true \
  -Dtriplea.port="40%i" \
  -Dtriplea.name="{{ bot_name }}%i_{{ bot_city }}" \
  -Dtriplea.lobby.uri="https://prod.triplea-game.org" \
  -Dtriplea.exit.on.game.end=true \
  -jar {{ bot_install_folder }}/bin/triplea-game-headless.jar
Restart=always

StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=bot%i


[Install]
WantedBy=multi-user.target
